# dingo/dingo.initd
#!/sbin/openrc-run

name="dingo"
command="/usr/bin/dingo"
command_args="$DINGO_OPTS"
command_background="yes"
pidfile="/var/run/dingo.pid"
